# src/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

# [program:irc]
# command=/usr/bin/sudo -u ircd /home/ircd/leaf/unreal start 

